# golink.service
[Unit]
Description=Go Links URL Shortener
After=network.target pihole-FTL.service
Wants=pihole-FTL.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/golink
ExecStart=/usr/local/bin/golink
Restart=always
RestartSec=5

# Environment configuration (optional, defaults work for most setups)
# Uncomment and modify if needed:
# Environment="PIHOLE_TARGET=http://127.0.0.1:8080"
# Environment="LISTEN_ADDR=:80"
# Environment="DB_FILE=go_links.json"

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
